# unparse/bin/crontab
# Crontab for undemcoracy

5 8 * * * $HOME/unparse/scraper/morningupdate

