# unparse/bin/crontab
# Crontab for undemcoracy

5 8 * * * $HOME/unparse-live/scraper/morningupdate

